# deploy/messagefeed-dev.service
[Unit]
Description=messageFeed production deployment
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/aroen/projects/Amoney/_Astu/go/go_st/Go_Pro/messageFeed
EnvironmentFile=/etc/messagefeed/messagefeed.env
ExecStart=/usr/local/bin/messageFeed-start cloudflare
ExecStop=/usr/local/bin/messageFeed-start stop
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
